# docker/test/core/import.ini
[program:import_product_live]
command=/usr/bin/php artisan import:product_live
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:import_product_offer_price_live]
command=/usr/bin/php artisan import:product_offer_prices live
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:import_product_offer_price_regular]
command=/usr/bin/php artisan import:product_offer_prices regular
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:import_product_offer_stocks]
command=/usr/bin/php artisan import:product_offer_stocks
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:import_products]
command=/usr/bin/php artisan import:products
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:import_pvz]
command=/usr/bin/php artisan import:pvz
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
