# docker/local/network_manager/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:monolith-ssh-tunnel]
command=/usr/local/bin/monolith-ssh-tunnel
autostart=true
autorestart=true
stdout_logfile=/dev/stdlog
stderr_logfile=/dev/stderr

